# roles/galaxy_labs_engine/files/labs_engine.socket
[Unit]
Description=Galaxy Labs Engine socket
[Socket]
ListenStream=/run/labs_engine.sock
[Install]
WantedBy=sockets.target
